# dashbot.service
[Unit]
Description=Runs dashbot

[Service]
Type=oneshot
ExecStart=/bin/ruby /bin/dashbot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
