# dashbot.service
[Unit]
Description=Runs dashbot

[Service]
ExecStart=/bin/rubyy /bin/dashbot

[Install]
WantedBy=multi-user.target
